# ChangerTaskTimer/files/kns_test.timer
[Unit]
description = hello
requires = kns.service

[Timer]
oncalendar = test
unit = kns.service

[Install]
wantedby = timers.targetOnCalendar

