# ChangerTaskTimer/files/kns_test2.timer
[Unit]
description=test2
requires=123

[Timer]
oncalendar=123
unit=123

[Install]
wantedby=123